# docs/launchd_templates/com.bot7.state-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bot7.state-snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>%BOT7_PATH%/.venv/bin/python</string>
        <string>scripts/state_snapshot_loop.py</string>
        <string>--interval-sec</string>
        <string>300</string>
    </array>
    <key>WorkingDirectory</key>
    <string>%BOT7_PATH%</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>%BOT7_PATH%/logs/launchd_state_snapshot.log</string>
    <key>StandardErrorPath</key>
    <string>%BOT7_PATH%/logs/launchd_state_snapshot.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:%BOT7_PATH%/.venv/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
