# docs/launchd_templates/com.bot7.refresh-ict.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bot7.refresh-ict</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/bot7/.venv/bin/python</string>
        <string>/Users/USERNAME/bot7/scripts/refresh_ict_levels.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/bot7</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>5</integer>
        <key>Minute</key><integer>0</integer>
        <key>Weekday</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/bot7/logs/launchd_ict.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/bot7/logs/launchd_ict.err</string>
</dict>
</plist>
